# archive_2025-07-04/redundant_files/configuration/soloheart.service
[Unit]
Description=SoloHeart Game Server with Monitor
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=%i
WorkingDirectory=/path/to/solo_heart
Environment=PATH=/path/to/venv/bin
ExecStart=/path/to/venv/bin/python launch_with_monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/solo_heart/logs /path/to/solo_heart/campaign_saves

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

[Install]
WantedBy=multi-user.target 